# karmada-scheduling-policy-observer/systemd/karmada-policy-observation.service
[Unit]
Description=Karmada production scheduling policy full observation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=karmada-observer
Group=karmada-observer
WorkingDirectory=/opt/karmada-scheduling-automation
EnvironmentFile=/etc/karmada-policy-observation.env
ExecStart=/opt/karmada-scheduling-automation/60-policy-observation.sh --once
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full
ReadWritePaths=/var/lib/karmada-policy-observation
